# etc/sv/carbon/run
#!/bin/bash

/opt/graphite/bin/carbon-cache.py start --nodaemon
