# config/sv/carbon/run
#!/bin/bash

rm -f /opt/graphite/storage/carbon-cache*.pid
/opt/graphite/bin/carbon-cache.py start --nodaemon
