# etc/sv/graphios/run
#!/bin/bash

su - $NAGIOS_USER -c '/usr/local/bin/graphios.py'
